# alt/phddns/phddns.conf
description   "Phddns daemon for this client"
author        "谢继雷 (Lenik) <lenik@bodz.net>"

start on (net-device-up IFACE!=lo)
stop on runlevel[!2345]

script

    sudo phddns

end script
