# lenik/etc/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.gkr18daddw installed on Thu Oct 22 11:22:06 2009)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
*/5 * * * * checkip
*/5 * * * * wget -qO- --http-user=lenik --http-password=bUbcHIaw "http://www.3322.org/dyndns/update?system=dyndns&hostname=lenik.3322.org" >>/tmp/3322.log

