# uni/xjl-coolbash/profile.d/color-ps.sh
PS1="\[\e]0;\u@\h: \w\a\]\${debian_chroot:+(\$debian_chroot)}"
PS1="$PS1[1;42m\u@\h:[0;47m\w[m
\$ "
