# uni/xjl-coolbash/profile.d/Makefile.am

sysconfdir = /etc
profilesdir = ${sysconfdir}/profile.d

dist_profiles_DATA = \
    alt-default.sh \
    large-history.sh \
    paths.sh
