# uni/xjl-coolbash/profile.d/paths.sh
export PATH=$PATH:~/bin:~/sbin
export SHLIB=/usr/lib/bash
