# lenik/etc/env.d/eci-secca/etc/profile.d/large-history.sh

HISTFILESIZE=100000000
HISTSIZE=1000000
