# lenik/etc/env.d/eci-secca/etc/profile.d/lenik.sh
#!/bin/bash

alias d='ls -FANl --color=auto'

function hello() {
	echo Hello, this is Lenik.
}

export EDITOR=vi
export PATH=$PATH:~/bin:~/sbin
