# lenik/etc/env.d/eci-secca/etc/profile.d/lenik.sh
#!/bin/bash

alias d='ls -FANo --color'

function hello() {
	echo Hello, this is Lenik.
}

export EDITOR=vi
export PATH=~/bin:$PATH
