# lenik/etc/env.d/eci-secca/etc/profile.d/su-keephist.sh
alias su='sudo su --login --preserve-environment'
