# lenik/etc/env.d/eci-secca/etc/profile.d/su-postgres.sh
alias pg='su postgres --login --preserve-environment'
