# uni/xjl-optalts/profile.d/optenv.sh
export PATH="/opt/bin:/opt/wine/bin:$PATH"
