# uni/toolenh/truecrypt-enh/mount-tc.conf

description "Mount TrueCrypt LAM Volumes"
author "Lenik <xjl@99jsj.com>"

start on local-filesystems
stop  on runlevel [06]

emits tc-mounted
emits tc-mapped

console output

exec mountall.tc
