# uni/xjl-sshproxy/sshproxy.conf.in
description    "Cool SOCKS5 thru *.sshproxy.bodz.net"
author        "谢继雷 (Lenik) <lenik@bodz.net>"

start on (net-device-up IFACE!=lo)
stop on runlevel[!2345]

script

    # -T disable pseudo-tty allocation
    # -f go to background after login but before command exec
    # -n stdin from /dev/null, must be used when ssh is run in bg.
    # -N no command
    # -D "dynamic" app-level port forwarding.
    sudo -usshproxy ssh -qTfn -ND *:7878 sshproxy.bodz.net

end script
