# core/etc/profile.d/99welcome.bat
rem #Lapiota successfully started up!

    rem mode con cols=67 lines=10
    echo.
    echo.
    echo.
    echo      ::: :::           ::
    echo       :   :             :
    echo       :   :             :
    echo       :   :   :::::     :     :::::   :::::  ::: :    :::::
    echo       : : :  :     :    :    :     : :     :  : : :  :     :
    echo       : : :  :::::::    :    :       :     :  : : :  :::::::
    echo       : : :  :          :    :       :     :  : : :  :
    echo        : :   :     :    :    :     : :     :  : : :  :     :
    echo        : :    :::::   :::::   :::::   :::::  :: : ::  :::::
    echo.
    echo.
    echo.

    rem wait for 3 seconds.
    readkey -t 3
