# core/etc/profile.d/zh_CN.sh

if false; then
    export LANG=zh_CN.GBK
    export LC_ALL=zh_CN.GBK
    export LC_CTYPE=zh_CN.GBK
    export OUTPUT_CHARSET=GBK
fi

export LESSCHARSET=latin1

alias dir='dir -N --color'
alias ls='ls --show-control-chars'

# export XMODIFIERS='@im=Chinput'
# e\stty cs8-istrip
# stty pass8
