# core/etc/cygwin/profile.d/zh_CN.sh

# alias ls='ls --'
# alias dir='dir -N --color'

export LANG=zh_CN.GBK
export OUTPUT_CHARSET=GBK
