# core/etc/cygwin/profile.d/zh_CN.sh

export LANG=zh_CN.GBK
export LC_ALL=zh_CN.GBK
export LC_CTYPE=zh_CN.GBK

# export XMODIFIERS='@im=Chinput'
# e\stty cs8-istrip
# stty pass8

export OUTPUT_CHARSET=GBK

export LESSCHARSET=latin1

alias dir='dir -N --color'
alias ls='ls --show-control-chars'
