# uni/toolenh/vim-enh/profile.d/coolvim.sh
export EDITOR=vi

# Disable <C-S>
stty stop ''
